# mw_te0820/project-spec/meta-mathworks/recipes-apps/mw-fs-overlay/files/common/fs-overlay/lib/systemd/system/usb_network.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/usb_network
Before=network.service
After=network_scripts.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/sbin/usb_network start
ExecStop=/usr/sbin/usb_network stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
