# mw_te0820/project-spec/meta-mathworks/recipes-apps/mw-fs-overlay/files/common/fs-overlay/lib/systemd/system/user_init.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/user_init
Before=user_app.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/sbin/user_init start
ExecStop=/usr/sbin/user_init stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
