# mw_zcu216/project-spec/meta-user/recipes-apps/mw-fs-overlay/files/services/backupSSHKeys.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/backupSSHKeys
Requires=sdcard_mount.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=no
ExecStart=/usr/sbin/backupSSHKeys start
ExecStop=/usr/sbin/backupSSHKeys stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
