# mw_zcu216/project-spec/meta-user/recipes-apps/mw-fs-overlay/files/services/inetd.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/inetd
After=network.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/sbin/inetd start
ExecStop=/usr/sbin/inetd stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
