# mw_zcu216/project-spec/meta-user/recipes-apps/mw-fs-overlay/files/services/network_scripts.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/network_scripts
Before=network.service
After=sdcard_mount.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/sbin/network_scripts start
ExecStop=/usr/sbin/network_scripts stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
