# mw_zcu216/project-spec/meta-user/recipes-apps/mw-fs-overlay/files/services/start_only.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/usr/sbin/start_only.sh

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/sbin/start_only.sh start
ExecStop=/usr/sbin/start_only.sh stop
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
